# V3.9/netfiremon-web.service
[Unit]
Description=Net-Fire-Monitor v3.0 – Web-Interface (unprivilegiert)
After=network-online.target netfiremon.service
Wants=network-online.target
# Web-Interface startet erst nachdem der Monitor-Prozess gestartet ist
Requires=netfiremon.service

# ── SICHERHEITS-ARCHITEKTUR (Gemini-Audit Fix 3) ─────────────────────────
# Dieser Prozess läuft als "netfiremon" (kein root, keine Capabilities).
# Er liest den State aus /opt/netfiremon/*.json und stellt das Web-Interface
# bereit. Keine Paketerfassung, kein iptables – keine root-Rechte nötig.
# ─────────────────────────────────────────────────────────────────────────

[Service]
Type=simple
User=netfiremon
Group=netfiremon
WorkingDirectory=/opt/netfiremon

# Gunicorn startet den Web-Server als unprivilegierter User
ExecStart=/opt/netfiremon/.venv/bin/gunicorn \
  --bind 0.0.0.0:5443 \
  --workers 2 \
  --timeout 30 \
  --certfile /opt/netfiremon/certs/cert.pem \
  --keyfile  /opt/netfiremon/certs/key.pem \
  --access-logfile - \
  --error-logfile - \
  "netfiremon_web:create_wsgi_app()"

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# ── Strikte Isolation (kein root, keine Capabilities) ─────────────────────
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true

# Nur die Verzeichnisse freigeben die der Web-Prozess wirklich braucht
# Nur data/ darf beschrieben werden – Code-Verzeichnis ist read-only
ReadWritePaths=/opt/netfiremon/data
ReadOnlyPaths=/opt/netfiremon /opt/netfiremon/certs

# Netzwerk: nur ausgehend auf Port 5443 binden (kein raw socket!)
# CapabilityBoundingSet bewusst leer – keinerlei erhöhte Rechte
CapabilityBoundingSet=

LimitNOFILE=65536
MemoryMax=256M

[Install]
WantedBy=multi-user.target
